# deploy/systemd/hrguru-team-evening-reminder.service
[Unit]
Description=Send HR Guru team evening WhatsApp reminder
Wants=network-online.target
After=network-online.target hrguru-whatsapp-relay.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/hrguru-ats
EnvironmentFile=-/home/ubuntu/hrguru-ats/.env
ExecStart=/home/ubuntu/hrguru-ats/venv/bin/python /home/ubuntu/hrguru-ats/scripts/send_team_whatsapp_reminder.py evening --send
