# deploy/systemd/hrguru-zero-submission-check.service
[Unit]
Description=Send HR Guru ATS 4 PM zero-submission check
Wants=network-online.target
After=network-online.target hrguru-whatsapp-relay.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/hrguru-ats
EnvironmentFile=-/home/ubuntu/hrguru-ats/.env
ExecStart=/home/ubuntu/hrguru-ats/venv/bin/python /home/ubuntu/hrguru-ats/scripts/no_submission_so_far_whatsapp_report.py --send
