# deploy/systemd/hrguru-weekly-submission-report.service
[Unit]
Description=Send HR Guru ATS weekly submission report
Wants=network-online.target
After=network-online.target hrguru-whatsapp-relay.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/hrguru-ats
EnvironmentFile=-/home/ubuntu/hrguru-ats/.env
ExecStart=/home/ubuntu/hrguru-ats/venv/bin/python /home/ubuntu/hrguru-ats/scripts/weekly_team_submission_whatsapp_report.py --send
